# docker/supervisor-makecert.conf
[program:makecert]
command=/bin/bash -c "/opt/angerona/makeCert.sh"
