# docker/supervisor-angerona.conf
[program:angerona]
command=/bin/bash -c "pserve /opt/angerona/production.ini"
