# docker/supervisor-makecert.conf
[program:makecert]
command=/bin/bash -c "/opt/angerona/makeCert.sh"
autorestart=unexpected
exitcodes=0
retrysecs=0
